# systemd/cloudflare-ddns.service
[Unit]
Description=Cloudflare Dynamic DNS Update
Documentation=https://github.com/tfmalt/utility-scripts
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment="CF_API_TOKEN_FILE=%h/.config/cloudflare/credentials"
# Add Homebrew to PATH (Linux Homebrew location)
Environment="PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin"

# Run the DDNS update script with explicit bash and sourcing
ExecStart=/bin/bash -c 'set -a; source ${CF_API_TOKEN_FILE}; set +a; exec %h/src/tfmalt/utility-scripts/scripts/cloudflare-ddns.sh'

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflare-ddns

[Install]
WantedBy=multi-user.target
