# systemd/cloudflare-ddns.service
[Unit]
Description=Cloudflare Dynamic DNS Update
Documentation=https://github.com/tfmalt/utility-scripts
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# Run the DDNS update script
# Note: Script will source credentials internally
ExecStart=/bin/bash -c 'source %h/.config/cloudflare/credentials && %h/src/tfmalt/utility-scripts/scripts/cloudflare-ddns.sh'

# Basic security hardening (compatible with user services)
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflare-ddns

[Install]
WantedBy=multi-user.target
